# dockerized-cron/crontab
*/1 * * * * bash /date-script.sh
*/1 * * * * echo 'crontab triggered'
# Single container cron example
